# deb/fio/DEBIAN/postinst
#!/bin/sh
# postinst script for blah
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see https://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in
    configure)
        ln -s /usr/local/bin/fio-nodeos /usr/local/bin/nodeos >/dev/null 2>&1 || true
        ln -s /usr/local/bin/clio /usr/local/bin/fio-cleos >/dev/null 2>&1 || true
        ln -s /usr/local/bin/clio /usr/local/bin/cleos >/dev/null 2>&1 || true
        ln -s /usr/local/bin/fio-keosd /usr/local/bin/keosd >/dev/null 2>&1 || true
	ln -s /etc/fio/nodeos/genesis-mainnet.json /etc/fio/nodeos/genesis.json >/dev/null 2>&1 || true
	[ ! -f /etc/fio/nodeos/config.ini ] && cp /etc/fio/nodeos/mainnet-config.ini /etc/fio/nodeos/config.ini
	[ ! -f /etc/fio/nodeos/logging.json ] && cp /etc/fio/nodeos/default-logging.json /etc/fio/nodeos/logging.json

	mkdir -p /etc/fio/nodeos/protocol_features
	chmod 0750 /etc/fio/nodeos/protocol_features
	chown -R fio:fio /etc/fio/nodeos/protocol_features
	mkdir -p /var/lib/fio/history /var/lib/fio/history_index
	chmod 0750 /var/lib/fio/history /var/lib/fio/history_index
	chown -R fio:fio /var/lib/fio/history /var/lib/fio/history_index
	systemctl daemon-reload || true

	echo
	echo "***************************************************************************"
	echo "*                                                                         *"
	echo "* Install complete, see /usr/local/share/fio/README for more information. *"
	echo "*  nodeos service is installed, but not started, it can be started with:  *"
	echo "*                                                                         *"
	echo "*               systemctl enable fio-nodeos.service                       *"
	echo "*               systemctl start fio-nodeos.service                        *"
	echo "*                                                                         *"
	echo "***************************************************************************"
	echo

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
